# docker-compose/docker/host-cron/crontab
* * * * * root sh /send_email.sh
#Single container cron example